# h.conf
[program:h]
command=python3 /home/curator2/h/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/h/test.err.log
stdout_logfile=/home/curator2/h/test.out.log